# templates/supervisord.ini

[program:mig-scheduler]
command={{ mig_src }}/bin/linux/amd64/mig-scheduler
directory={{ mig_src }}
user={{ mig_user }}
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
redirect_stderr=true
stdout_logfile=/var/log/supervisor/mig-scheduler.log

[program:mig-api]
command={{ mig_src }}/bin/linux/amd64/mig-api
directory={{ mig_src }}
user={{ mig_user }}
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/mig-api.log

